# graveboards.service
[Unit]
Description=Graveboards Backend Service
Requires=docker.service network-online.target
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="POSTGRES_PASSWORD=${POSTGRES_PASSWORD}"
Environment="SESSION_SECRET=${SESSION_SECRET}"

# User and group (optional - create dedicated user)
# User=graveboards
# Group=graveboards

# Security options
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

# Working directory (update this to your actual deployment path)
WorkingDirectory=/path/to/graveboards-deploy

# Start command (update path to your deployment directory)
ExecStart=/usr/bin/docker-compose -f /path/to/graveboards-deploy/docker-compose.prod.yml up -d

# Stop command
ExecStop=/usr/bin/docker-compose -f /path/to/graveboards-deploy/docker-compose.prod.yml down

# Restart command
ExecReload=/usr/bin/docker-compose -f /path/to/graveboards-deploy/docker-compose.prod.yml restart

# Restart policy
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=graveboards

[Install]
WantedBy=multi-user.target
